# contrib/systemd/status-herald-curtain.service
[Unit]
Description=status-herald per-tab curtain focus agent (dispatcher: poll or event stream)
After=network-online.target

[Service]
Type=simple
# The systemd user PATH is not your shell PATH: it has no nvm/asdf/volta, so
# `herald` (and the `node` its shebang needs) will not be found. Either symlink
# both into ~/.local/bin, or uncomment and adjust the line below.
#Environment=PATH=%h/.nvm/versions/node/v22.12.0/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/usr/bin/env herald curtain arm-all
ExecStart=/usr/bin/env bash %h/.local/share/status-herald/run.sh
# SIGTERM kills the poll loop before its trailing reveal-all can run. Without
# this, stopping the service leaves every covered session behind its card.
ExecStopPost=/usr/bin/env herald curtain reveal-all
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
